# deploy/systemd.service
[Unit]
Description=Warpsmith Battle Simulator
Documentation=https://github.com/fixemer90-stack/Warpsmith
After=network.target

[Service]
Type=simple
User=warpsmith
WorkingDirectory=/opt/warpsmith
Environment=HOSTING=true
Environment=DB_PATH=/var/lib/warpsmith/simulator.db
Environment=PORT=8000
ExecStart=/opt/warpsmith/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
